# deploy/photo-backend.service
[Unit]
Description=Photo Battle Backend
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=photo-battle
WorkingDirectory=__INSTALL_DIR__/backend

ExecStart=__INSTALL_DIR__/backend/target/release/photo-battle-backend
Restart=always
RestartSec=5

# 从项目 .env 加载敏感变量
EnvironmentFile=__INSTALL_DIR__/.env

# 非敏感路径变量
Environment=BASE_URL=http://localhost:8080
Environment=DATABASE_PATH=__INSTALL_DIR__/data/photos.db
Environment=UPLOAD_DIR=__INSTALL_DIR__/uploads
Environment=ARTIMUSE_ENABLED=false
Environment=VOLCENGINE_ENABLED=true
Environment=VOLCENGINE_URL=http://127.0.0.1:8001

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
