# deploy/volcengine.service
[Unit]
Description=VolcEngine Aesthetic Scoring Service
After=network.target

[Service]
Type=simple
User=photo-battle
WorkingDirectory=__INSTALL_DIR__/backend/scripts

ExecStart=__INSTALL_DIR__/.venv/bin/python volcengine_aesthetic_server.py
Restart=always
RestartSec=10

# 从项目 .env 加载 HF_TOKEN 等
EnvironmentFile=__INSTALL_DIR__/.env

Environment=VOLCENGINE_PORT=8001
Environment=VOLCENGINE_DEVICE=cpu

# 首次启动需下载 CLIP 模型（~1.7GB），超时放宽
TimeoutStartSec=600

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
